# etc/cron.d/updater
PATH=/usr/local/bin:/usr/bin:${PATH}

0 3 * * * pi cd /home/pi/jlock && bash scripts/updater.sh
